# deploy/speedfog-racing.service
[Unit]
Description=SpeedFog Racing API
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=speedfog
Group=speedfog
WorkingDirectory=/opt/speedfog-racing/server
EnvironmentFile=/opt/speedfog-racing/.env
ExecStart=/opt/speedfog-racing/server/.venv/bin/uvicorn speedfog_racing.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/data/SpeedFog
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
